# data/vnc.socket
[Unit]
Before=vnc.service

[Socket]
ListenStream=5900
Accept=no

[Install]
WantedBy=sockets.target
